# assets/supervisord/postfix.conf
[program:postfix]
startsecs = 0
autostart = true
autorestart = false
command = postfix -c /etc/postfix start
process_name = master
directory = /etc/postfix
